# systemd/pi-tv-kiosk.service.example
[Unit]
Description=PI TV Kiosk Browser
After=network-online.target docker.service graphical.target
Wants=network-online.target docker.service

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/pi-tv-docker
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/user/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
Environment=KIOSK_URL=http://localhost:3000
Environment=AUDIO_STACK=pipewire
Environment=KIOSK_PROFILE_DIR=/tmp/kiosk-chromium-profile
Environment=KIOSK_FOREGROUND=1
ExecStart=/home/user/pi-tv-docker/scripts/restart-kiosk-browser.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
